# contrib/init/babycoind.conf
description "Babycoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BABYCOIND_BIN="/usr/bin/babycoind"
env BABYCOIND_USER="babycoin"
env BABYCOIND_GROUP="babycoin"
env BABYCOIND_PIDDIR="/var/run/babycoind"
# upstart can't handle variables constructed with other variables
env BABYCOIND_PIDFILE="/var/run/babycoind/babycoind.pid"
env BABYCOIND_CONFIGFILE="/etc/babycoin/babycoin.conf"
env BABYCOIND_DATADIR="/var/lib/babycoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # babycoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BABYCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run babycoind."
        echo "The setting must appear in $BABYCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Babycoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BABYCOIND_PIDDIR"
    chmod 0755 "$BABYCOIND_PIDDIR"
    chown $BABYCOIND_USER:$BABYCOIND_GROUP "$BABYCOIND_PIDDIR"
    chown $BABYCOIND_USER:$BABYCOIND_GROUP "$BABYCOIND_CONFIGFILE"
    chmod 0660 "$BABYCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BABYCOIND_PIDFILE" \
    --chuid $BABYCOIND_USER:$BABYCOIND_GROUP \
    --exec "$BABYCOIND_BIN" \
    -- \
    -pid="$BABYCOIND_PIDFILE" \
    -conf="$BABYCOIND_CONFIGFILE" \
    -datadir="$BABYCOIND_DATADIR" \
    -disablewallet \
    -daemon

